# scripts/systemd/homebridge-unicorn-hat.service
[Unit]
Description=Homebridge Unicorn Hat Daemon

[Service]
ExecStart=/usr/local/bin/homebridge-unicorn-hat
Environment="HUH_LOG_LEVEL=debug"

[Install]
WantedBy=multi-user.target
